# .config/systemd/user/dns-dangling-monitor.timer
[Unit]
Description=Run DNS dangling monitor daily

[Timer]
OnCalendar=*-*-* 04:30:00 Asia/Tokyo
RandomizedDelaySec=20m
Persistent=true
Unit=dns-dangling-monitor.service

[Install]
WantedBy=timers.target
